# monitoring-test.service
[Unit]
Description=Monitor process "test" and check HTTPS endpoint
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitoring-test.sh
User=root